# backend.service
[Unit]
Description=BBIFather Backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/bbifather/backend
Environment=TELEGRAM_BOT_TOKEN=YOUR_BOT_TOKEN
Environment=TELEGRAM_CHAT_ID=YOUR_CHAT_ID
ExecStart=/usr/bin/python3 -m gunicorn -w 2 -b 127.0.0.1:8000 main:app
Restart=always

[Install]
WantedBy=multi-user.target